# automations/launchd/com.neo.easyoutfit_codex_bridge.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>com.neo.easyoutfit_codex_bridge</string>

    <key>ProgramArguments</key>
    <array>
      <string>/bin/bash</string>
      <string>/Users/neo/Desktop/closetgptrenew/scripts/run_local_codex_bridge_easyoutfit.sh</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/neo/Desktop/closetgptrenew</string>

    <key>EnvironmentVariables</key>
    <dict>
      <key>PATH</key>
      <string>/opt/homebrew/bin:/opt/homebrew/sbin:/usr/local/bin:/usr/local/sbin:/usr/bin:/bin:/usr/sbin:/sbin</string>
      <key>EASYOUTFIT_LOCAL_CODEX_WORKSPACE_SLUG</key>
      <string>easyoutfitapp</string>
      <key>EASYOUTFIT_LOCAL_CODEX_WORKER_SUFFIX</key>
      <string>easyoutfit-codex-bridge</string>
      <key>EASYOUTFIT_LOCAL_CODEX_BRIDGE_ENV_FILE</key>
      <string>/Users/neo/Desktop/closetgptrenew/.env.easyoutfit_codex_bridge</string>
    </dict>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>ThrottleInterval</key>
    <integer>10</integer>

    <key>ProcessType</key>
    <string>Background</string>

    <key>StandardOutPath</key>
    <string>/Users/neo/Desktop/closetgptrenew/logs/easyoutfit_codex_bridge.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/neo/Desktop/closetgptrenew/logs/easyoutfit_codex_bridge.log</string>
  </dict>
</plist>
